# debian/postinst
#!/bin/bash
#
# postinst script for linuxmuster-base
# thomas@linuxmuster.net
# 20191204
# GPL v3
#

# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package
#

# get constants
source /usr/share/linuxmuster/defaults.sh

case "$1" in

 configure)

  # add lmn paths to python environment
  usersite="$(python3 -c 'import site; site._script()' --user-site)"
  userpth="$usersite/lmn.pth"
  mkdir -p "$usersite"
  echo "$LIBDIR" > "$userpth"
  echo "$SETUPDIR" >> "$userpth"

  # clean and bytecompile modules
  find "$LIBDIR" -name \*.pyc -exec rm '{}' \;
  find "$LIBDIR" -type d -name __pycache__ -exec rm -r '{}' \; 2> /dev/null || true
  python3 -m compileall "$LIBDIR"

  # create ssl-cert group
  groupadd --force --system ssl-cert

  # permissions
  chgrp ssl-cert "$SSLDIR" -R
  chmod 750 "$SSLDIR"
  chmod 700 "$SECRETDIR"
  sysvol="/var/lib/samba/sysvol"
  [ -d "$sysvol" ] && find "$sysvol" -type d -exec chmod 775 '{}' \;

  # create firewall fullchain pem file
  fwcert="$(echo "$FWFULLCHAIN" | sed -e 's|.fullchain.|.cert.|')"
  if [ ! -s "$FWFULLCHAIN" -a -s "$fwcert" ]; then
    echo "Creating $(basename "$FWFULLCHAIN")."
    cat "$fwcert" "$CACERT" > "$FWFULLCHAIN"
    chmod 640 "$FWFULLCHAIN"
  fi

  # provide cacert.pem for clients if not present
  if [ -n "$domainname" -a -s "$CACERT" ]; then
    sysvoltlsdir="$(echo "$SYSVOLTLSDIR" | sed -e 's|@@domainname@@|'"$domainname"'|')"
    sysvolpemfile="$sysvoltlsdir/$(basename "$CACERT")"
    [ -d "$sysvoltlsdir" ] || mkdir -p "$sysvoltlsdir"
    if [ -d "$sysvoltlsdir" -a ! -e "$sysvolpemfile" ]; then
      echo "Providing $sysvolpemfile."
      cp "$CACERT" "$sysvolpemfile"
    fi
  fi

  # timesyncd update, add firewall ip as parent ntp server
  conf='/etc/systemd/timesyncd.conf'
  search="NTP=${firewallip}"
  if ! grep -q ^"$search" "$conf"; then
    echo "Updating time server related timesyncd configuration."
    cp "$conf" "$conf".dpkg-bak
    echo "$search" >> "$conf"
    systemctl daemon-reload
  fi

  # smb.conf update, enable time server
  tserver="$(testparm -l -s --parameter-name="time server" 2>/dev/null | tr A-Z a-z)"
  if [ "$tserver" != "yes" ];then
    echo "Updating time server related samba configuration."
    conf="/etc/samba/smb.conf"
    cp "$conf" "$conf".dpkg-bak
    linuxmuster-modini -i "$conf" -s global -o 'time server' -v Yes -r samba-ad-dc
  fi

  # dhcp update, server ip as ntp server
  conf='/etc/dhcp/dhcpd.conf'
  search="option ntp-servers ${firewallip}"
  if grep -q ^"$search" "$conf"; then
    echo "Updating time server related dhcp configuration."
    replace="${search/$firewallip/$serverip}"
    cp "$conf" "$conf".dpkg-bak
    sed -i "s|^$search|$replace|" "$conf"
    systemctl restart isc-dhcp-server.service
  fi

  # enable ntp service, change firewall name to ip
  if ! systemctl status ntp.service &> /dev/null; then
    conf="/etc/ntp.conf"
    echo "Enabling ntp service."
    cp "$conf" "$conf".dpkg-bak
    sed -i "s|^pool firewall|pool $firewallip|" "$conf"
    systemctl enable ntp.service
    systemctl start ntp.service
  fi

 ;;

 abort-upgrade|abort-remove|abort-deconfigure)
 ;;

 *)
  echo "postinst called with unknown argument \`$1'" >&2
  exit 1
 ;;

esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
